# backend/crontab
# Backup every day at 2:00 AM
0 2 * * *  /app/backup.sh >> /var/log/cron.log 2>&1
